# youtube_downloader/rootfs/etc/services.d/youtube-downloader/run
#!/usr/bin/with-contenv bashio
set -e

record_ytdlp_update() {
    /venv/bin/python - "$1" <<'PY'
from __future__ import annotations

import json
import os
import sys
from datetime import UTC, datetime
from pathlib import Path

state_file = Path("/data/jobs/ytdlp_update.json")
state_file.parent.mkdir(parents=True, exist_ok=True)
try:
    state = json.loads(state_file.read_text(encoding="utf-8"))
    if not isinstance(state, dict):
        state = {}
except (OSError, json.JSONDecodeError):
    state = {}
now = datetime.now(UTC).isoformat()
state["last_attempt"] = now
if sys.argv[1] == "success":
    state["last_success"] = now
    state.pop("last_error", None)
else:
    state["last_error"] = "Aktualizacja yt-dlp przy starcie nie powiodła się."
temporary = state_file.with_suffix(".tmp")
temporary.write_text(json.dumps(state, ensure_ascii=False, indent=2), encoding="utf-8")
os.replace(temporary, state_file)
PY
}

bashio::log.info "Aktualizuję yt-dlp do najnowszej wersji przed uruchomieniem aplikacji..."
if /venv/bin/python -m pip install --no-cache-dir --retries 3 --timeout 20 --upgrade yt-dlp; then
    record_ytdlp_update success
else
    bashio::log.warning "Aktualizacja yt-dlp nie powiodła się. Uruchamiam z obecną wersją."
    record_ytdlp_update failure
fi

EXTERNAL_PORT="$(bashio::config 'external_port')"
BIND_ARGS=(--bind 0.0.0.0:8099)

if bashio::config.true 'allow_external_port'; then
    if [ "${EXTERNAL_PORT}" = "8099" ]; then
        bashio::log.info "Dostęp zewnętrzny używa głównego portu 8099."
    else
        BIND_ARGS+=(--bind "0.0.0.0:${EXTERNAL_PORT}")
        bashio::log.info "Włączono publiczny dostęp bez Ingress na porcie ${EXTERNAL_PORT}."
        if [ "${EXTERNAL_PORT}" != "999" ]; then
            bashio::log.warning "Jeśli port ${EXTERNAL_PORT} ma działać z sieci, ustaw zgodne mapowanie w sekcji Sieć dodatku."
        fi
    fi
else
    bashio::log.info "Publiczny dostęp bez Ingress jest wyłączony."
fi

bashio::log.info "Uruchamiam Media Web Downloader na porcie 8099."
exec /venv/bin/gunicorn \
    --chdir / \
    "${BIND_ARGS[@]}" \
    --workers 1 \
    --threads 4 \
    --timeout 0 \
    --access-logfile - \
    --error-logfile - \
    "app.main:app"
